# custom-packages/luci-app-unicast-proxy/files/etc/init.d/unicast_proxy
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=1

CONFIG=unicast_proxy

# ── 防火墙规则管理 ────────────────────────────────────
# 用固定命名的 section（unicast_proxy_fw_<id>）便于精确增删
add_firewall_rule() {
	local cfg="$1" listen_port="$2" name="$3"
	local sec="unicast_proxy_fw_${cfg}"

	uci -q delete firewall.$sec 2>/dev/null
	uci set firewall.$sec='rule'
	uci set firewall.$sec.name="UnicastProxy-${name:-$cfg}"
	uci set firewall.$sec.src='wan'
	uci set firewall.$sec.proto='tcp'
	uci set firewall.$sec.dest_port="$listen_port"
	uci set firewall.$sec.target='ACCEPT'
}

remove_firewall_rule() {
	local cfg="$1"
	uci -q delete firewall.unicast_proxy_fw_${cfg} 2>/dev/null
}

# ── 单条规则处理 ──────────────────────────────────────
start_rule() {
	local cfg="$1"
	local enable name upstream_host upstream_port listen_port wan_access

	config_get_bool enable     "$cfg" enable 0
	config_get name            "$cfg" name "$cfg"
	config_get upstream_host   "$cfg" upstream_host
	config_get upstream_port   "$cfg" upstream_port
	config_get listen_port     "$cfg" listen_port
	config_get_bool wan_access "$cfg" wan_access 0

	# 每次重启都先清掉本规则旧的防火墙条目，再按当前配置决定是否重新添加
	remove_firewall_rule "$cfg"

	[ "$enable" -eq 1 ] || return 0
	[ -n "$upstream_host" ] || return 0
	[ -n "$upstream_port" ] || return 0
	[ -n "$listen_port" ]   || return 0

	procd_open_instance "proxy_${cfg}"
	procd_set_param command /usr/bin/socat \
		TCP-LISTEN:${listen_port},fork,reuseaddr,nodelay,keepalive \
		TCP:${upstream_host}:${upstream_port},keepalive
	procd_set_param respawn
	procd_set_param stderr 1
	procd_close_instance

	if [ "$wan_access" -eq 1 ]; then
		add_firewall_rule "$cfg" "$listen_port" "$name"
	fi
}

stop_rule_firewall() {
	remove_firewall_rule "$1"
}

start_service() {
	config_load "$CONFIG"
	config_foreach start_rule rule

	uci commit firewall
	/etc/init.d/firewall reload >/dev/null 2>&1
}

stop_service() {
	config_load "$CONFIG"
	config_foreach stop_rule_firewall rule

	uci commit firewall
	/etc/init.d/firewall reload >/dev/null 2>&1
}

reload_service() {
	stop
	start
}

service_triggers() {
	procd_add_reload_trigger "$CONFIG"
}
